# supervisor.conf
[program:nginx]
# nginx switches user on its own
#user=www-data
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
redirect_stderr=true
autorestart=true
